# uwsgi.conf.sample
description "uWSGI Emperor Starter"
start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /bin/bash /home/talaikis/uwsgi.sh
